# bluez/code/dbus-1.9.4/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/disk4/liangshen/practice/bluez_outdir/var/run/dbus/system_bus_socket
